# agent/tcf-agent_1.7.0-r0_armhf/DEBIAN/postinst
#!/bin/bash

if [ ! -e /etc/systemd/system/tcf-agent.service ]
then
	systemctl enable tcf-agent.service
	systemctl start tcf-agent.service
fi


